# warpengine/files/etc/init.d/warpengine
#!/bin/sh /etc/rc.common

# OpenWrt procd init script for WarpEngine (2025 极简兼容版)
# 优化项：TCP BBR, CAKE QoS, FullCone NAT, 智能 DNS

START=95
STOP=15

USE_PROCD=1

# 使用纯 Shell 逻辑比较延迟（无需 bc）
fastest_dns() {
    local dns_list="119.29.29.29 223.5.5.5 180.76.76.76 114.114.114.114"
    local best_dns="119.29.29.29"
    local min_int="9999999" # 初始设为一个很大的整数

    for dns in $dns_list; do
        # 获取 ping 结果，超时 1 秒
        local line=$(ping -c 1 -W 1 "$dns" 2>/dev/null | grep 'time=')
        if [ -n "$line" ]; then
            # 提取延迟数字
            local ms_raw=$(echo "$line" | sed 's/.*time=\([0-9.]*\).*/\1/')
            # 方案 1 核心：删除小数点，将 12.345 转换为 12345 整数
            local ms_int=$(echo "$ms_raw" | tr -d '.')
            
            if [ -n "$ms_int" ]; then
                # 使用 Shell 原生整数比较 -lt (less than)
                if [ "$ms_int" -lt "$min_int" ]; then
                    min_int="$ms_int"
                    best_dns="$dns"
                fi
            fi
        fi
    done
    echo "$best_dns"
}

start_service() {
    config_load 'warpengine'

    local enabled sw_offload hw_offload ipv6_offload mss_clamping game_process bbr fullcone dns_acc dns_mode dns_server
    config_get_bool enabled main enabled 0
    [ "$enabled" -eq 1 ] || return 0

    # 使用 oneshot 模式，适用于配置下发类脚本
    procd_open_instance
    procd_set_param command /bin/true
    procd_set_param type oneshot

    # 1. 网络加速与防火墙 (支持 fw4/nftables)
    config_get_bool sw_offload main sw_flow_offload 0
    config_get_bool hw_offload main hw_flow_offload 0
    config_get_bool ipv6_offload main ipv6_offload 0
    config_get_bool mss_clamping main mss_clamping 1
    config_get_bool fullcone main fullcone_nat 0

    uci set firewall.@defaults[-1].flow_offloading="$sw_offload"
    uci set firewall.@defaults[-1].hw_flow_offloading="$hw_offload"
    uci set firewall.@defaults[-1].flow_offloading_6="$ipv6_offload"
    uci set firewall.@defaults[-1].tcp_mss_clamping="$mss_clamping"
    uci set firewall.@defaults[-1].fullcone="$fullcone"
    uci commit firewall
    /etc/init.d/firewall reload
    logger -t "WarpEngine" "Firewall 加速已应用 (FullCone: $fullcone)"

    # 2. TCP 拥塞控制 (BBR)
    config_get_bool bbr main bbr 0
    if [ "$bbr" -eq 1 ]; then
        # 尝试加载内核模块，如果失败则跳过
        sysctl -w net.core.default_qdisc=fq >/dev/null 2>&1
        sysctl -w net.ipv4.tcp_congestion_control=bbr >/dev/null 2>&1
        logger -t "WarpEngine" "TCP BBR 已启用"
    else
        sysctl -w net.ipv4.tcp_congestion_control=cubic >/dev/null 2>&1
        sysctl -w net.core.default_qdisc=fq_codel >/dev/null 2>&1
    fi

    # 3. 游戏优化 (CAKE 队列管理)
    config_get_bool game_process main game_process 0
    # 自动获取 LAN 设备名，2025 年多为 br-lan
    local lan_if=$(uci get network.lan.device 2>/dev/null || echo "br-lan")
    tc qdisc del dev "$lan_if" root 2>/dev/null
    if [ "$game_process" -eq 1 ]; then
        # 确保系统中安装了 kmod-sched-cake
        tc qdisc add dev "$lan_if" root cake 2>/dev/null
        logger -t "WarpEngine" "已在 $lan_if 部署 CAKE 调度"
    fi

    # 4. DNS 加速设置
    config_get_bool dns_acc main dns_acc 0
    if [ "$dns_acc" -eq 1 ]; then
        config_get dns_mode main dns_mode 'static'
        local target_dns
        if [ "$dns_mode" = "auto" ]; then
            logger -t "WarpEngine" "正在通过整数对比法寻找最快 DNS..."
            target_dns=$(fastest_dns)
        else
            config_get target_dns main dns_server '119.29.29.29'
        fi
        
        uci -q del_list dhcp.@dnsmasq[0].server
        uci add_list dhcp.@dnsmasq[0].server="$target_dns"
        uci set dhcp.@dnsmasq[0].noresolv='1'
        logger -t "WarpEngine" "DNS 设置为: $target_dns"
    else
        uci set dhcp.@dnsmasq[0].noresolv='0'
        uci -q del_list dhcp.@dnsmasq[0].server
    fi
    
    uci commit dhcp
    /etc/init.d/dnsmasq restart
    
    procd_close_instance
}

stop_service() {
    # 恢复防火墙
    uci set firewall.@defaults[-1].flow_offloading='0'
    uci set firewall.@defaults[-1].hw_flow_offloading='0'
    uci set firewall.@defaults[-1].flow_offloading_6='0'
    uci set firewall.@defaults[-1].tcp_mss_clamping='1' 
    uci set firewall.@defaults[-1].fullcone='0'
    uci commit firewall
    /etc/init.d/firewall reload

    # 清除 CAKE
    local lan_if=$(uci get network.lan.device 2>/dev/null || echo "br-lan")
    tc qdisc del dev "$lan_if" root 2>/dev/null

    # 恢复 DNS
    uci set dhcp.@dnsmasq[0].noresolv='0'
    uci -q del_list dhcp.@dnsmasq[0].server
    uci commit dhcp
    /etc/init.d/dnsmasq restart
    
    # 恢复 TCP 默认
    sysctl -w net.ipv4.tcp_congestion_control=cubic >/dev/null 2>&1
    sysctl -w net.core.default_qdisc=fq_codel >/dev/null 2>&1
    
    logger -t "WarpEngine" "WarpEngine 已停止，配置已还原。"
}

