# warpengine/files/etc/init.d/warpengine
#!/bin/sh /etc/rc.common

START=95
USE_PROCD=1

# DNS 优选算法
fastest_dns() {
    local dns_list="119.29.29.29 223.5.5.5 180.76.76.76 114.114.114.114"
    local best_dns="119.29.29.29"
    local min_int="9999999"
    for dns in $dns_list; do
        local line=$(ping -c 1 -W 1 "$dns" 2>/dev/null | grep 'time=')
        if [ -n "$line" ]; then
            local ms_int=$(echo "$line" | sed 's/.*time=\([0-9.]*\).*/\1/' | tr -d '.')
            if [ "$ms_int" -lt "$min_int" ]; then min_int="$ms_int"; best_dns="$dns"; fi
        fi
    done
    echo "$best_dns"
}

start_service() {
    config_load 'warpengine'
    local master
    config_get_bool master main enabled 0

    if [ "$master" -eq 0 ]; then
        logger -t "WarpEngine" "总控开关已关闭，还原每个节点出厂值..."
        stop_service
        return 0
    fi

    logger -t "WarpEngine" "总控开关已启用，开始检查各功能节点状态..."

    # 1. 软件/硬件/IPv6 加速/MSS/FullCone (防火墙组)
    local sw hw v6 mss fc
    config_get_bool sw main sw_flow_offload 0
    config_get_bool hw main hw_flow_offload 0
    config_get_bool v6 main ipv6_offload 0
    config_get_bool mss main mss_clamping 1
    config_get_bool fc main fullcone_nat 0

    uci set firewall.@defaults[-1].flow_offloading="$sw"
    uci set firewall.@defaults[-1].hw_flow_offloading="$hw"
    uci set firewall.@defaults[-1].flow_offloading_6="$v6"
    uci set firewall.@defaults[-1].tcp_mss_clamping="$mss"
    uci set firewall.@defaults[-1].fullcone="$fc"
    uci commit firewall
    /etc/init.d/firewall reload
    
    [ "$sw" -eq 1 ] && logger -t "WarpEngine" "软件加速：已启用" || logger -t "WarpEngine" "软件加速：已禁用 (还原出厂)"
    [ "$hw" -eq 1 ] && logger -t "WarpEngine" "硬件加速：已启用" || logger -t "WarpEngine" "硬件加速：已禁用 (还原出厂)"
    [ "$v6" -eq 1 ] && logger -t "WarpEngine" "全栈加速：已启用" || logger -t "WarpEngine" "全栈加速：已禁用 (还原出厂)"
    [ "$fc" -eq 1 ] && logger -t "WarpEngine" "FullCone NAT：已启用" || logger -t "WarpEngine" "FullCone NAT：已禁用 (还原出厂)"

    # 2. TCP BBR
    local bbr
    config_get_bool bbr main bbr 0
    if [ "$bbr" -eq 1 ]; then
        sysctl -w net.core.default_qdisc=fq >/dev/null 2>&1
        sysctl -w net.ipv4.tcp_congestion_control=bbr >/dev/null 2>&1
        logger -t "WarpEngine" "TCP BBR 算法：已启用"
    else
        sysctl -w net.ipv4.tcp_congestion_control=cubic >/dev/null 2>&1
        sysctl -w net.core.default_qdisc=fq_codel >/dev/null 2>&1
        logger -t "WarpEngine" "TCP BBR 算法：已禁用 (还原出厂)"
    fi

    # 3. CAKE
    local cake
    config_get_bool cake main game_process 0
    local lan_if=$(uci -q get network.lan.device || echo "br-lan")
    tc qdisc del dev "$lan_if" root 2>/dev/null
    if [ "$cake" -eq 1 ]; then
        tc qdisc add dev "$lan_if" root cake 2>/dev/null
        logger -t "WarpEngine" "游戏转发 (CAKE)：已启用"
    else
        logger -t "WarpEngine" "游戏转发 (CAKE)：已禁用 (还原出厂)"
    fi

    # 4. DNS
    local dns_acc
    config_get_bool dns_acc main dns_acc 0
    if [ "$dns_acc" -eq 1 ]; then
        local mode target_dns
        config_get mode main dns_mode 'static'
        [ "$mode" = "auto" ] && target_dns=$(fastest_dns) || config_get target_dns main dns_server '119.29.29.29'
        uci -q del dhcp.@dnsmasq.server
        uci add_list dhcp.@dnsmasq.server="$target_dns"
        uci set dhcp.@dnsmasq.noresolv='1'
        logger -t "WarpEngine" "智能 DNS 加速：已启用 ($target_dns)"
    else
        uci set dhcp.@dnsmasq.noresolv='0'
        uci -q del_list dhcp.@dnsmasq.server
        logger -t "WarpEngine" "智能 DNS 加速：已禁用 (还原出厂)"
    fi
    uci commit dhcp
    /etc/init.d/dnsmasq restart
}

stop_service() {
    # 强制所有节点回到出厂值逻辑
    uci set firewall.@defaults[-1].flow_offloading='0'
    uci set firewall.@defaults[-1].hw_flow_offloading='0'
    uci set firewall.@defaults[-1].flow_offloading_6='0'
    uci set firewall.@defaults[-1].fullcone='0'
    uci commit firewall && /etc/init.d/firewall reload

    local lan_if=$(uci -q get network.lan.device || echo "br-lan")
    tc qdisc del dev "$lan_if" root 2>/dev/null

    sysctl -w net.ipv4.tcp_congestion_control=cubic >/dev/null 2>&1
    sysctl -w net.core.default_qdisc=fq_codel >/dev/null 2>&1

    uci set dhcp.@dnsmasq.noresolv='0'
    uci -q del_list dhcp.@dnsmasq.server
    uci commit dhcp && /etc/init.d/dnsmasq restart
}

