# warpengine/files/etc/init.d/warpengine
#!/bin/sh /etc/rc.common

START=95
USE_PROCD=1

log() {
	logger -t "warpengine" "$*"
}

fastest_dns() {
	local dns_list="119.29.29.29 223.5.5.5 180.76.76.76 114.114.114.114"
	local best_dns="119.29.29.29"
	local min_ms=9999
	for dns in $dns_list; do
		local output=$(dig @"$dns" example.com +short +time=1 +tries=1 2>/dev/null)
		if [ -n "$output" ]; then
			local query_time=$(dig @"$dns" example.com +stats +time=1 +tries=1 2>/dev/null | grep "Query time:" | awk '{print $4}')
			if [ -n "$query_time" ] && [ "$query_time" -lt "$min_ms" ]; then
				min_ms="$query_time"
				best_dns="$dns"
			fi
		fi
	done
	echo "$best_dns"
}

get_flowtable_devices() {
	local devs=""
	local lan_dev=$(uci -q get network.lan.device)
	local wan_dev=$(uci -q get network.wan.device)
	[ -n "$lan_dev" ] && devs="$lan_dev"
	[ -n "$wan_dev" ] && devs="$devs $wan_dev"
	if [ -n "$devs" ]; then
		echo "{ $(echo $devs | tr ' ' ', ') }"
	else
		echo "{ }"
	fi
}

check_dependencies() {
	[ -x /usr/sbin/nft ] || { log "ERROR: nftables not available"; return 1; }
	[ -x /usr/bin/dig ] || { log "ERROR: bind-dig not installed"; return 1; }
	[ -x /sbin/tc ] || { log "WARN: tc not available (CAKE will be skipped)"; }
	return 0
}

start_service() {
	check_dependencies || return 1

	config_load 'warpengine'
	local master; config_get_bool master main enabled 0
	[ "$master" -eq 0 ] && { stop_service; return 0; }

	local sw hw v6 mss fc mwan3_compat
	config_get_bool sw main sw_flow_offload 0
	config_get_bool hw main hw_flow_offload 0
	config_get_bool v6 main ipv6_offload 0
	config_get_bool mss main mss_clamping 1
	config_get_bool fc main fullcone_nat 0
	config_get_bool mwan3_compat main mwan3_compat 0

	uci set firewall.@defaults[-1].flow_offloading="$sw"
	uci set firewall.@defaults[-1].hw_flow_offloading="$hw"
	uci set firewall.@defaults[-1].flow_offloading_6="$v6"
	uci set firewall.@defaults[-1].tcp_mss_clamping="$mss"
	uci set firewall.@defaults[-1].fullcone="$fc"

	local ft_priority=0; [ "$mwan3_compat" -eq 1 ] && ft_priority=1500
	local ft_devices=$(get_flowtable_devices)
	nft delete flowtable inet fw4 ft 2>/dev/null
	nft add flowtable inet fw4 ft "{ hook ingress priority $ft_priority; devices = $ft_devices; }" 2>/dev/null

	uci commit firewall
	/etc/init.d/firewall reload

	local bbr cake dns_acc packet_steering irq_balance
	config_get_bool bbr main bbr 0
	config_get_bool cake main game_process 0
	config_get_bool packet_steering main packet_steering 0
	config_get_bool irq_balance main irq_balance 0

	if [ "$bbr" -eq 1 ]; then
		sysctl -w net.core.default_qdisc=fq >/dev/null 2>&1
		sysctl -w net.ipv4.tcp_congestion_control=bbr >/dev/null 2>&1
	else
		sysctl -w net.ipv4.tcp_congestion_control=cubic >/dev/null 2>&1
		sysctl -w net.core.default_qdisc=fq_codel >/dev/null 2>&1
	fi

	local lan_if=$(uci -q get network.lan.device || echo "br-lan")
	tc qdisc del dev "$lan_if" root 2>/dev/null
	if [ "$cake" -eq 1 ] && tc qdisc help 2>/dev/null | grep -q cake; then
		tc qdisc add dev "$lan_if" root cake 2>/dev/null
	fi

	uci set network.globals.packet_steering="$packet_steering"
	uci commit network

	if [ "$irq_balance" -eq 1 ] && [ -f "/etc/init.d/irqbalance" ]; then
		/etc/init.d/irqbalance enable
		/etc/init.d/irqbalance start
	elif [ -f "/etc/init.d/irqbalance" ]; then
		/etc/init.d/irqbalance disable
		/etc/init.d/irqbalance stop
	fi

	config_get_bool dns_acc main dns_acc 0
	if [ "$dns_acc" -eq 1 ]; then
		local mode target_dns
		config_get mode main dns_mode 'static'
		if [ "$mode" = "auto" ]; then
			target_dns=$(fastest_dns)
		else
			config_get target_dns main dns_server '119.29.29.29'
		fi
		uci -q del dhcp.@dnsmasq.server
		uci add_list dhcp.@dnsmasq.server="$target_dns"
		uci set dhcp.@dnsmasq.noresolv='1'
	else
		uci set dhcp.@dnsmasq.noresolv='0'
		uci -q del_list dhcp.@dnsmasq.server
	fi
	uci commit dhcp
	/etc/init.d/dnsmasq restart

	log "WarpEngine started successfully."
}

stop_service() {
	uci set firewall.@defaults[-1].flow_offloading='0'
	uci set firewall.@defaults[-1].hw_flow_offloading='0'
	uci set firewall.@defaults[-1].flow_offloading_6='0'
	uci set firewall.@defaults[-1].fullcone='0'
	uci commit firewall

	local ft_devices=$(get_flowtable_devices)
	nft delete flowtable inet fw4 ft 2>/dev/null
	nft add flowtable inet fw4 ft "{ hook ingress priority 0; devices = $ft_devices; }" 2>/dev/null
	/etc/init.d/firewall reload

	uci set network.globals.packet_steering='0'
	uci commit network

	if [ -f "/etc/init.d/irqbalance" ]; then
		/etc/init.d/irqbalance disable
		/etc/init.d/irqbalance stop
	fi

	local lan_if=$(uci -q get network.lan.device || echo "br-lan")
	tc qdisc del dev "$lan_if" root 2>/dev/null

	sysctl -w net.ipv4.tcp_congestion_control=cubic >/dev/null 2>&1
	sysctl -w net.core.default_qdisc=fq_codel >/dev/null 2>&1

	uci set dhcp.@dnsmasq.noresolv='0'
	uci -q del_list dhcp.@dnsmasq.server
	uci commit dhcp
	/etc/init.d/dnsmasq restart

	log "WarpEngine stopped."
}

reload_service() {
	log "Reloading WarpEngine..."
	stop_service
	start_service
}
